# i3/systemd/user/gnome-background-change.service
[Unit]
Description=Rotate GNOME background
 
[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/bash /home/queljin/bin/rotate_bg.sh
