# i3/systemd/user/i3-background-change.service
[Unit]
Description=Rotate i3 background
 
[Service]
Type=oneshot
ExecStart=/usr/bin/bash /home/queljin/bin/rotate_bg_feh.sh
